# cpumoder.service
[Unit]
Description=CPU go brrrrrrrrrrrrrrrrrrr
[Service]
ExecStart=/bin/cpumod
[Install]
WantedBy=basic.target
