# cpusdg.service
[Unit]
Description=cpusdg
[Service]
ExecStart=/usr/bin/cpusdg
[Install]
WantedBy=basic.target
