# cpusdg.service
[Unit]
Description=CPU go brrrrrrrrrrrrr
[Service]
ExecStart=/usr/bin/cpugobrr
[Install]
WantedBy=basic.target
